# init.d/proxy-firewall
#!/sbin/openrc-run
name="proxy-firewall"
description="Firewall for US-COX proxy"

depend() {
    need net
    after firewall
}

start() {
    /bin/sh /etc/proxy-whitelist/firewall.sh ssh-on
    /bin/sh /etc/proxy-whitelist/firewall.sh wl-apply
}

stop() {
    /bin/sh /etc/proxy-whitelist/firewall.sh wl-off
    /bin/sh /etc/proxy-whitelist/firewall.sh ssh-off
}
